# gunicorn.conf
[program:gunicorn]
command=gunicorn --workers=3 main_flask:app -b localhost:8001
directory=/code
